# docker/rc.local
#!/bin/sh
/home/app/yarn/bin/yarn.yml.sh > /home/app/yarn/config/yarn.yml
chown app:app /home/app/yarn/config/yarn.yml
crontab -r
crontab 2>/dev/null -l | { cat; echo "0  * * * * /home/app/yarn/bin/merit-ranks.sh";    } | crontab -
crontab 2>/dev/null -l | { cat; echo "0  0 * * * /home/app/yarn/bin/export-csv.sh";     } | crontab -
crontab 2>/dev/null -l | { cat; echo "0  5 * * * /home/app/yarn/bin/export-xml.sh";     } | crontab -
crontab 2>/dev/null -l | { cat; echo "0 10 * * * /home/app/yarn/bin/refresh-scores.sh"; } | crontab -
